# src/configurations/bluetooth/bluetooth-autoaccept.service
[Unit]
Description=Bluetooth Auto-Accept Agent Service
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/bluetooth-autoaccept.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
