# src/configurations/bluetooth/bt-auto-connect.service.tmpl
[Unit]
Description=Bluetooth Auto-Connect Service
After=bluetooth.service
Requires=bluetooth.service

[Service]
ExecStart=/usr/local/bin/bt-auto-connect
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
