# src/configurations/profile.d/xdg-runtime.sh
# export XDG runtime dir & D-Bus session for any login
if [ "$(id -u)" -ge 1000 ] && [ -d "/run/user/$(id -u)" ]; then
  export XDG_RUNTIME_DIR="/run/user/$(id -u)"
  export DBUS_SESSION_BUS_ADDRESS="unix:path=$XDG_RUNTIME_DIR/bus"
fi